# docs/braubar.conf
[program:braubar]
command=./braubar
directory=/home/oli/dev/braubar-pi/bin

